# deploy/backup/hauldr-postgres-backup.timer
[Unit]
Description=Daily Hauldr Postgres backup

[Timer]
OnCalendar=*-*-* 04:00:00 UTC
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
