# deploy/backup/hauldr-postgres-backup.service
[Unit]
Description=Hauldr Postgres off-site backup (pg_dumpall → object store)
After=docker.service
Wants=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/hauldr-postgres-backup.sh
